# MPICH/NOVNC Layer/conf.d/fluxbox.conf
[program:fluxbox]
command=startfluxbox
autorestart=true
