# docker/conf.d/jupyter.conf
[program:jupyternotebook]
command=/shft/app/iparaview-kernel/start_jupyter.sh
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
